# overrides/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:postgres]
command=/usr/bin/postmaster -D /etc/postgresql
user=postgres
